# nginx/systemd/nginx-watcher.path
[Unit]
Description=Nginx Service Reloader
PartOf=nginx.service

[Path]
Unit=nginx-watcher.service
PathChanged=/etc/dsiprouter/certs/

[Install]
WantedBy=multi-user.target